# timeshift-snapshots.path
[Unit]
Description=Monitors for new Timeshift snapshots

[Path]
PathModified=/run/timeshift/backup/timeshift-btrfs/snapshots
Unit=update-systemd-boot-snapshots.service

[Install]
WantedBy=multi-user.target

