# update-systemd-boot-snapshots.service
[Unit]
Description=Update systemd boot snapshots

[Service]
Type=oneshot
ExecStart=update-systemd-boot-snapshots